# system_files/shared/usr/lib/systemd/system/swtpm-workaround.service.d/10-order-after-rechunker-group-fix.conf
[Unit]
After=rechunker-group-fix.service
