# system_files/shared/usr/lib/systemd/system/aurora-zfs-deprecation-notifier.timer
[Unit]
Description=Run aurora-zfs-deprecation-notifier 20 mins after boot

[Timer]
OnBootSec=20min
Unit=aurora-zfs-deprecation-notifier.service

[Install]
WantedBy=timers.target
